# supervisord.ini
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:cde-server]
command=/cde-server